# configuration-as-code/stow/systemd/etc/systemd/system/conorganizer-sqlite-backup.service
[Unit]
Description=Back up Conorganizer SQLite database

[Service]
Type=oneshot
ExecStart=/usr/local/bin/conorganizer-sqlite-backup
