# configuration-as-code/stow/systemd/etc/systemd/system/conorganizer-images-backup.service
[Unit]
Description=Back up Conorganizer event images

[Service]
Type=oneshot
ExecStart=/usr/local/bin/conorganizer-images-backup
