# configuration-as-code/stow/systemd/etc/systemd/system/conorganizer-sqlite-backup.timer
[Unit]
Description=Run Conorganizer SQLite database backup

[Timer]
OnCalendar=*:0/15
Persistent=true

[Install]
WantedBy=timers.target
